# macOS/profile.d/path.sh
export PATH="/usr/local/sbin:/usr/local/bin:$PATH"

export PATH="/usr/local/opt/terraform@0.11/bin:$PATH"
