# macOS/profile.d/tokens.sh
export ATLAS_TOKEN=""
export VAGRANTCLOUD_TOKEN=""
